# boot/gentoo-freevo-mga
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need lircd
}

start() {
	ebegin "Starting freevo"

	export DISPLAY=""
        sysctl -w dev.rtc.max-user-freq=1024 > /dev/null
	/usr/bin/freevo 2>/dev/null >/dev/null&
	eend $?
}

stop() {
	ebegin "Stopping freevo"
	/usr/bin/freevo stop 2>/dev/null >/dev/null&
	eend 0
}
